# build/rootfs/var/lib/cache/fallback/opt/gearbox/GearboxOS/etc/upstart/wplib-box-start.conf
# GearBox - startup script.

description	"GearBox"

start on started docker
task
console output

env PATH="/opt/gearbox/bin:/opt/gearbox/sbin:/root/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

pre-start script
	if [ -x /opt/gearbox/bin/box ]
	then
		if [ ! -f /etc/box-scripts-updated ]
		then
			echo "GearBox: Running first time provisioning."
			/opt/gearbox/bin/box first-time-provision
		fi
	fi
end script

script
	if [ -x /opt/gearbox/bin/box ]
	then
		echo "GearBox: Starting up..."
		/opt/gearbox/bin/box startup
	else
		echo "GearBox: ERROR: W021 - box-scripts not present."
		echo "Pull down the latest box-scripts by doing the following:"
		echo "	sudo git clone -q https://github.com/gearboxworks/box-scripts/ /opt/gearbox"
		echo "Run an update:"
		echo "	/opt/gearbox/bin/box self-update"
	fi
end script

