# build/rootfs/etc/profile.d/box-profile.sh
/opt/gearbox/GearboxOS/etc/profile.d/box-profile.sh