# build/rootfs/etc/local.d/01-boot-check.stop
#!/bin/bash

echo ""
echo "################################################################################"
echo "# Stopping GearboxOS."
echo "# $(date)"

/etc/init.d/gearbox-check stop

box shutdown

/etc/init.d/docker-check stop
