# build/rootfs/var/lib/cache/fallback/opt/gearbox/GearboxOS/etc/crontabs/boxuser
0 * * * * /opt/gearbox/bin/box database backup
0 * * * * /opt/gearbox/bin/box component update
